# gen_service/supervisor/app.conf
[program:genapp]
command=/usr/bin/runapp.sh
autostart=true
autorestart=true
